# supervisor/supervisord.cfg
[supervisord]
loglevel=info
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stopsignal=QUIT

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
stopsignal=QUIT

[program:sidekiq]
command=su git -c /supervisor/start_sidekiq
autostart=true
autorestart=true
stopsignal=QUIT

[program:unicorn]
command=su git -c /supervisor/start_unicorn
autostart=true
autorestart=true
stopsignal=QUIT
